# project/crontab
55 23 * * * /var/lib/postgresql/data > /var/backups
