# project/crontab
55 23 * * * pg_dump /var/lib/postgresql/data > /var/backups
